# qr3_83.conf
[program:qr3_83]
command=python3 /home/curator3/qr3_83/main.py
autostart=true
autorestart=true
stderr_logfile=/home/curator3/qr3_83/test.err.log
stdout_logfile=/home/curator3/qr3_83/test.err.log
